# systemd/fresh-boot-heartbeat.service
[Unit]
Description=Boot freshness heartbeat - periodic timestamp update
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/touch /var/lib/.boot-freshness-last-seen
